# cookbooks/hello_service/templates/default/hello_service.init.conf.erb
#!upstart
#description "hello_service startup"

respawn

start on startup
stop on shutdown

setuid vagrant
setgid vagrant

script
  cd /vagrant/source/hello_service
  node app <%= node.hello_service.port_number %>
end script
